# target/linux/rockchip/armv8/base-files/etc/hotplug.d/net/40-net-smp-affinity
#!/bin/sh

[ "$ACTION" = add ] || exit

get_device_irq() {
	local device="$1"
	local line
	local seconds="0"

	# wait up to 10 seconds for the irq/device to appear
	while [ "${seconds}" -le 10 ]; do
		line=$(grep -E -m 1 "${device}\$" /proc/interrupts) && break
		seconds="$(( seconds + 2 ))"
		sleep 2
	done
	echo ${line} | sed 's/:.*//'
}

set_interface_core() {
	local core_mask="$1"
	local interface="$2"
	local device="$3"

	[ -z "${device}" ] && device="$interface"

	local irq=$(get_device_irq "$device")

	echo "${core_mask}" > /proc/irq/${irq}/smp_affinity
}

case "$(board_name)" in
ariaboard,photonicat|\
armsom,sige3|\
ezpro,mrkaio-m68s|\
firefly,rk3568-roc-pc|\
friendlyarm,nanopi-r5c|\
lunzn,fastrhino-r66s)
	set_interface_core 2 "eth0"
	set_interface_core 4 "eth1"
	;;
friendlyarm,nanopi-r2c|\
friendlyarm,nanopi-r2c-plus|\
friendlyarm,nanopi-r2s|\
radxa,rockpi-e|\
xunlong,orangepi-r1-plus|\
xunlong,orangepi-r1-plus-lts)
	set_interface_core 2 "eth0"
	set_interface_core 4 "eth1" "xhci-hcd:usb[0-9]+"
	;;
friendlyarm,nanopi-r4s|\
friendlyarm,nanopi-r4se|\
friendlyarm,nanopi-r4s-enterprise)
	set_interface_core 10 "eth0"
	set_interface_core 20 "eth1"
	;;
friendlyarm,nanopi-r5s)
	set_interface_core 2 "eth0"
	set_interface_core 4 "eth1"
	set_interface_core 8 "eth2"
	;;
hinlink,opc-h68k)
	set_interface_core "0-3" "eth0"
	set_interface_core "1" "eth1"
	set_interface_core "2" "eth2-0"
	set_interface_core "2" "eth2-16"
	set_interface_core "1" "eth2-18"
	set_interface_core "3" "eth3-0"
	set_interface_core "3" "eth3-18"
	set_interface_core "1" "eth3-16"
	rfc=32768
	sysctl net.core.rps_sock_flow_entries=$rfc
	for fileRfc in $(ls /sys/class/net/eth*/queues/rx-*/rps_flow_cnt)
	do
		eth_name=$(echo "$fileRfc"|awk -F/ '{print $5}')
		echo $rfc > $fileRfc
	done
	#0101
	echo 9 > /sys/class/net/eth0/queues/rx-0/rps_cpus
	echo 9 > /sys/class/net/eth1/queues/rx-0/rps_cpus
	#0011
	echo 3 > /sys/class/net/eth2/queues/rx-0/rps_cpus
	echo 3 > /sys/class/net/eth2/queues/rx-1/rps_cpus
	#1100
	echo c > /sys/class/net/eth3/queues/rx-0/rps_cpus
	echo c > /sys/class/net/eth3/queues/rx-1/rps_cpus
	/usr/sbin/ethtool -K eth0 tso on sg on tx on
	/usr/sbin/ethtool -K eth1 tso on sg on tx on
	/usr/sbin/ethtool -K eth2 tso on sg on tx on
	/usr/sbin/ethtool -K eth3 tso on sg on tx on
	;;
lunzn,fastrhino-r68s)
	set_interface_core 1 "eth0"
	set_interface_core 2 "eth1"
	set_interface_core 4 "eth2"
	set_interface_core 8 "eth3"
	;;
lyt,t68m)
	set_interface_core 1 "lan1"
	set_interface_core 2 "lan2"
	set_interface_core 4 "lan3"
	set_interface_core 8 "lan4"
	;;
esac
