# package/ctcgfw/luci-app-233ddns/root/etc/uci-defaults/luci_233ddns
#!/bin/sh

uci -q batch <<-EOF >/dev/null
	delete ucitrack.@233ddns[-1]
	add ucitrack 233ddns
	set ucitrack.@233ddns[-1].init=233ddns
	commit ucitrack
EOF

rm -f /tmp/luci-indexcache
exit 0