# package/ctcgfw/luci-theme-my/root/etc/uci-defaults/30_luci-theme-my
#!/bin/sh
uci batch <<-EOF
	set luci.themes.AdvancedTomatomy=/luci-static/my
	commit luci
EOF
exit 0
