# package/ctcgfw/luci-app-mwol/root/etc/uci-defaults/40_luci-mwol
#!/bin/sh

uci -q batch <<-EOF >/dev/null
	delete ucitrack.@mwol[-1]
	add ucitrack mwol
	set ucitrack.@mwol[-1].init=mwol
	commit ucitrack
EOF

rm -f /tmp/luci-indexcache
exit 0
