# target/linux/hisilicon/base-files/etc/init.d/bootcount
#!/bin/sh /etc/rc.common

# 启动优先级设为 99，即在所有服务之后运行
START=99

boot() {
    # 海思盒子暂时不需要重置计数器
    return 0
}