# target/linux/qualcommax/base-files/etc/uci-defaults/999_auto-restart.sh
#!/bin/sh

/etc/init.d/network restart
/etc/init.d/odhcpd restart
/etc/init.d/rpcd restart

exit 0
