# package/ctcgfw/luci-theme-infinityfreedom/root/etc/uci-defaults/10_luci-theme-infinityfreedom
#!/bin/sh
uci batch <<-EOF
	set luci.themes.infinityfreedom=/luci-static/infinityfreedom
	set luci.main.mediaurlbase=/luci-static/infinityfreedom
	commit luci
EOF
exit 0
