# target/linux/qualcommax/ipq60xx/base-files/etc/hotplug.d/firmware/11-ath11k-caldata
#!/bin/sh

[ -e /lib/firmware/$FIRMWARE ] && exit 0

. /lib/functions/caldata.sh

board=$(board_name)

case "$FIRMWARE" in
"ath11k/IPQ6018/hw1.0/cal-ahb-c000000.wifi.bin")
	case "$board" in
	8devices,mango-dvk|\
	alfa-network,ap120c-ax)
		caldata_extract "0:ART" 0x1000 0x20000
		;;
	anysafe,e1|\
	cambiumnetworks,xe3-4)
		caldata_extract "0:ART" 0x1000 0x10000
		;;
	cmiot,ax18|\
	glinet,gl-ax1800|\
	glinet,gl-axt1800|\
	linksys,mr7350|\
	linksys,mr7500|\
	netgear,wax214|\
	netgear,wax610|\
	netgear,wax610y|\
	qihoo,360v6|\
	redmi,ax5|\
	xiaomi,ax1800|\
	zn,m2)
		caldata_extract "0:art" 0x1000 0x10000
		;;
	tplink,eap610od|\
	tplink,eap623od-hd-v1|\
	tplink,eap625od-hd-v1)
		caldata_from_file "/tmp/factory_data/radio" 0 0x10000
		;;
	yuncore,fap650)
		caldata_extract "0:art" 0x1000 0x20000
		;;
	jdcloud,re-ss-01|\
	jdcloud,re-cs-02|\
	link,nn6000-v1|\
	link,nn6000-v2|\
	redmi,ax5-jdcloud)
		caldata_extract_mmc "0:ART" 0x1000 0x10000
		;;
	esac
	;;
"ath11k/QCN9074/hw1.0/cal-pci-0000:01:00.0.bin")
	case "$board" in
	anysafe,e1|\
	cambiumnetworks,xe3-4)
		caldata_extract "0:ART" 0x26800 0x20000
		;;
	linksys,mr7500)
		caldata_extract "0:art" 0x26800 0x20000
		;;
	jdcloud,re-cs-02)
		caldata_extract_mmc "0:ART" 0x26800 0x20000
		;;
	esac
	;;
*)
	exit 1
	;;
esac
