# target/linux/ramips/base-files/etc/init.d/set-irq-affinity
#!/bin/sh /etc/rc.common
#MTK-Wi-Fi

START=99

start() {
	/sbin/mtk-wifi.sh &
}
