# target/linux/ipq807x/base-files/etc/init.d/set-irq-affinity
#!/bin/sh /etc/rc.common

START=99

start() {
	for eth in /sys/class/net/lan* /sys/class/net/wan /sys/class/net/10g-*; do
		[ -d "$eth" ] || continue
		for i in "$eth"/queues/rx-*/rps_cpus; do echo "f" > "$i"; done
		for i in "$eth"/queues/rx-*/rps_flow_cnt; do echo "4096" > "$i"; done
		for i in "$eth"/queues/tx-*/xps_cpus; do echo "f" > "$i"; done
	done

	echo "32768" > "/proc/sys/net/core/rps_sock_flow_entries"
}
