# target/linux/qualcommax/ipq807x/base-files/etc/init.d/bootcount
#!/bin/sh /etc/rc.common

START=99

boot() {
	case $(board_name) in
	edgecore,eap102|\
	yuncore,ax880)
		fw_setenv upgrade_available 0
		# Unset changed flag after sysupgrade complete
		fw_setenv changed
	;;
	linksys,mx4200v1|\
	linksys,mx4200v2|\
	linksys,mx4300|\
	linksys,mx5300|\
	linksys,mx8500)
		mtd resetbc s_env || true
	;;
	netgear,wax620|\
	netgear,wax630)
		fw_setenv boot_count 0
	;;
	redmi,ax6-stock|\
	xiaomi,ax3600-stock|\
	xiaomi,ax9000-stock)
		# OTA handling should not be used. Reset it just in case.
		fw_setenv flag_ota_reboot 0
		# Not strictly needed but useful to handle partition crash condition
		fw_setenv flag_boot_success 1
	;;
	esac
}
