# target/linux/ath79/mikrotik/base-files/etc/hotplug.d/firmware/10-ath9k-eeprom
#!/bin/sh

[ -e /lib/firmware/$FIRMWARE ] && exit 0

. /lib/functions/caldata.sh
. /lib/functions/mikrotik-caldata.sh

ath9k_eeprom_die() {
	echo "ath9k eeprom: " "$*"
	exit 1
}

ath9k_patch_fw_mac() {
	local mac=$1
	local mac_offset=$(($2))
	local chksum_offset=$3
	local xor_mac
	local xor_fw_mac
	local xor_fw_chksum

	[ -z "$mac" -o -z "$mac_offset" ] && return

	[ -n "$chksum_offset" ] && {
		xor_mac=${mac//:/}
		xor_mac="${xor_mac:0:4} ${xor_mac:4:4} ${xor_mac:8:4}"

		xor_fw_mac=$(hexdump -v -n 6 -s $mac_offset -e '/1 "%02x"' /lib/firmware/$FIRMWARE)
		xor_fw_mac="${xor_fw_mac:0:4} ${xor_fw_mac:4:4} ${xor_fw_mac:8:4}"

		xor_fw_chksum=$(hexdump -v -n 2 -s $chksum_offset -e '/1 "%02x"' /lib/firmware/$FIRMWARE)
		xor_fw_chksum=$(xor $xor_fw_chksum $xor_fw_mac $xor_mac)

		printf "%b" "\x${xor_fw_chksum:0:2}\x${xor_fw_chksum:2:2}" | \
			dd of=/lib/firmware/$FIRMWARE conv=notrunc bs=1 seek=$(($chksum_offset)) count=2
	}

	macaddr_2bin $mac | dd of=/lib/firmware/$FIRMWARE conv=notrunc oflag=seek_bytes bs=6 seek=$mac_offset count=1
}

board=$(board_name)

case "$FIRMWARE" in
"ath9k-eeprom-ahb-18100000.wmac.bin")
	case $board in
	mikrotik,routerboard-lhg-2nd|\
	mikrotik,routerboard-sxt-5nd-r2)
		mikrotik_caldata_extract "art" 0x1000 0x440
		ath9k_patch_fw_mac $(macaddr_add "$mac_base" +1)
		;;
	mikrotik,routerboard-wap-g-5hact2hnd)
		mikrotik_caldata_extract "art" 0x1000 0x440
		ath9k_patch_fw_mac $(macaddr_add $(mtd_get_mac_binary art 0x10) +2)
		;;
	*)
		ath9k_eeprom_die "board $board is not supported yet"
		;;
	esac
	;;
esac
