# package/ctcgfw/luci-app-vssr/root/etc/uci-defaults/luci-vssr
#!/bin/sh

uci -q batch <<-EOF >/dev/null
	delete ucitrack.@vssr[-1]
	add ucitrack vssr
	set ucitrack.@vssr[-1].init=vssr
	commit ucitrack
	delete firewall.vssr
	set firewall.vssr=include
	set firewall.vssr.type=script
	set firewall.vssr.path=/var/etc/vssr.include
	set firewall.vssr.reload=1
	commit firewall
EOF

/etc/init.d/dnscrypt-proxy stop
/etc/init.d/dnscrypt-proxy disable
/etc/init.d/privoxy stop
/etc/init.d/privoxy disable
killall -q -9 v2ray-plugin

/usr/share/vssr/gfw2ipset.sh
rm -f /tmp/luci-indexcache
exit 0
