# package/ctcgfw/luci-theme-btmod/root/etc/uci-defaults/luci-theme-btmod
#!/bin/sh
uci batch <<-EOF
	set luci.themes.Btmod=/luci-static/btmod
	set luci.main.mediaurlbase=/luci-static/btmod
	commit luci
EOF
exit 0
