# package/lienol/luci-theme-bootstrap-mod/root/etc/uci-defaults/30_luci-theme-bootstrap-mod
#!/bin/sh
uci batch <<-EOF
	set luci.themes.Bootstrap_Mod=/luci-static/bootstrap_mod
	commit luci
EOF
exit 0
