# package/ctcgfw/open-app-filter/files/appfilter.init
#!/bin/sh /etc/rc.common

START=96
start() {
	lsmod | grep -q oaf 2>/dev/null || modprobe oaf
	gen_class.sh /etc/appfilter/feature.cfg
	appfilter.sh
}

stop() {
	echo "stop appfilter"
	rmmod oaf >/dev/null 2>&1
}
