# package/ntlf9t/luci-app-rtorrent/src/etc/init.d/rtorrent
#!/bin/sh /etc/rc.common

START=99
STOP=99

start() {
  HOME=/root screen -dmS rtorrent nice -19 rtorrent
}

boot() {
  start
}

stop() {
  killall rtorrent
}

