# package/system/fstools/files/fstab.init
#!/bin/sh /etc/rc.common
# SPDX-License-Identifier: GPL-2.0-only
#
# Copyright (C) 2013-2020 OpenWrt.org

START=11

boot() {
	/sbin/block mount
}

start() {
	/sbin/block mount
}

restart() {
	/sbin/block umount
	/sbin/block mount
}

restart() {
	start
}

stop() {
	/sbin/block umount
}
