# package/rooter/ext-rooter-basic/files/etc/init.d/bannerset
#!/bin/sh /etc/rc.common
# Copyright (C) 2013 OpenWrt.org

START=80
USE_PROCD=1
