# target/linux/mediatek/filogic/base-files/etc/uci-defaults/05_fix-compat-version
. /lib/functions.sh

case "$(board_name)" in
	bananapi,bpi-r3|\
	bananapi,bpi-r3-mini)
	uci set system.@system[0].compat_version="1.2"
	uci commit system
	;;
esac

exit 0
