# target/linux/rockchip/armv8/base-files/etc/init.d/phy-leds
#!/bin/sh /etc/rc.common

START=30

led_set() {
	local path="/sys/class/leds/$1"
	local params="$2"
	local value="$3"

	[ -d "$path" ] || return 1

	for param in $params; do
		echo "$value" > "$path/$param"
	done
}

boot() {
	case "$(board_name)" in
	lunzn,fastrhino-r68s|\
	ezpro,mrkaio-m68s)
		led_set "stmmac-0:01:amber:lan" "trigger" "netdev"
		led_set "stmmac-0:01:amber:lan" "device_name" "eth0"
		led_set "stmmac-0:01:amber:lan" "link_10 link_100 link_1000" "0"
		led_set "stmmac-0:01:amber:lan" "rx tx" "1"
		led_set "stmmac-0:01:green:lan" "trigger" "netdev"
		led_set "stmmac-0:01:green:lan" "device_name" "eth0"
		led_set "stmmac-0:01:green:lan" "link_10 link_100 link_1000" "1"
		led_set "stmmac-0:01:green:lan" "rx tx" "0"
		led_set "stmmac-1:01:amber:lan" "trigger" "netdev"
		led_set "stmmac-1:01:amber:lan" "device_name" "eth1"
		led_set "stmmac-1:01:amber:lan" "link_10 link_100 link_1000" "0"
		led_set "stmmac-1:01:amber:lan" "rx tx" "1"
		led_set "stmmac-1:01:green:lan" "trigger" "netdev"
		led_set "stmmac-1:01:green:lan" "device_name" "eth1"
		led_set "stmmac-1:01:green:lan" "link_10 link_100 link_1000" "1"
		led_set "stmmac-1:01:green:lan" "rx tx" "0"
		;;
	friendlyarm,nanopi-r3s)
		led_set "enp1s0-1::lan" "link_10 link_100 link_1000" "1"
		led_set "enp1s0-1::lan" "rx tx" "0"
		led_set "enp1s0-2::lan" "link_10 link_100 link_1000" "0"
		led_set "enp1s0-2::lan" "rx tx" "1"
		led_set "stmmac-0:01:amber:wan" "trigger" "netdev"
		led_set "stmmac-0:01:amber:wan" "device_name" "eth0"
		led_set "stmmac-0:01:amber:wan" "link_10 link_100 link_1000" "0"
		led_set "stmmac-0:01:amber:wan" "rx tx" "1"
		led_set "stmmac-0:01:green:wan" "trigger" "netdev"
		led_set "stmmac-0:01:green:wan" "device_name" "eth0"
		led_set "stmmac-0:01:green:wan" "link_10 link_100 link_1000" "1"
		led_set "stmmac-0:01:green:wan" "rx tx" "0"
		;;
	friendlyarm,nanopi-r4s)
		led_set "enp1s0-1::lan" "link_10 link_100 link_1000" "1"
		led_set "enp1s0-1::lan" "rx tx" "0"
		led_set "enp1s0-2::lan" "link_10 link_100 link_1000" "0"
		led_set "enp1s0-2::lan" "rx tx" "1"
		;;
	friendlyarm,nanopi-r5s)
		led_set "stmmac-0:01:amber:lan" "trigger" "netdev"
		led_set "stmmac-0:01:amber:lan" "device_name" "eth0"
		led_set "stmmac-0:01:amber:lan" "link_10 link_100 link_1000" "1"
		led_set "stmmac-0:01:amber:lan" "rx tx" "0"
		led_set "stmmac-0:01:green:lan" "trigger" "netdev"
		led_set "stmmac-0:01:green:lan" "device_name" "eth0"
		led_set "stmmac-0:01:green:lan" "link_10 link_100 link_1000" "0"
		led_set "stmmac-0:01:green:lan" "rx tx" "1"
		;;
	huake,guangmiao-g4c)
                led_set "enp1s0-0::lan" "link_10 link_100 link_1000" "1"
                led_set "enp1s0-0::lan" "rx tx" "0"
                led_set "enp1s0-1::lan" "link_10 link_100 link_1000" "0"
                led_set "enp1s0-1::lan" "rx tx" "1"
		;;
	lyt,t68m)
		led_set "stmmac-0:01:amber:lan" "trigger" "netdev"
		led_set "stmmac-0:01:amber:lan" "device_name" "lan1"
		led_set "stmmac-0:01:amber:lan" "link_10 link_100 link_1000" "0"
		led_set "stmmac-0:01:amber:lan" "rx tx" "1"
		led_set "stmmac-0:01:green:lan" "trigger" "netdev"
		led_set "stmmac-0:01:green:lan" "device_name" "lan1"
		led_set "stmmac-0:01:green:lan" "link_10 link_100 link_1000" "1"
		led_set "stmmac-0:01:green:lan" "rx tx" "0"
		led_set "stmmac-1:01:amber:lan" "trigger" "netdev"
		led_set "stmmac-1:01:amber:lan" "device_name" "lan2"
		led_set "stmmac-1:01:amber:lan" "link_10 link_100 link_1000" "0"
		led_set "stmmac-1:01:amber:lan" "rx tx" "1"
		led_set "stmmac-1:01:green:lan" "trigger" "netdev"
		led_set "stmmac-1:01:green:lan" "device_name" "lan2"
		led_set "stmmac-1:01:green:lan" "link_10 link_100 link_1000" "1"
		led_set "stmmac-1:01:green:lan" "rx tx" "0"
		;;
	esac
}
