# target/linux/rockchip/armv8/base-files/etc/init.d/reload-sdio-wifi
#!/bin/sh /etc/rc.common

START=99

boot() {
	case "$(board_name)" in
	ariaboard,photonicat)
		( sleep 15s; wifi; ) &
		;;
	armsom,sige3|\
	firefly,rk3568-roc-pc)
		( sleep 60s; wifi; ) &
		;;
	esac
}
