# package/jsda/luci-app-rtorrent/src/etc/uci-defaults/rtorrent
#!/bin/sh
/etc/init.d/rtorrent enable
/etc/init.d/rtorrent start
rm -fr /tmp/luci-indexcache /tmp/luci-modulecache
exit 0

